# run/systemd/generator/-.mount
# Automatically generated by systemd-fstab-generator

[Unit]
SourcePath=/etc/fstab
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
Before=local-fs.target

[Mount]
What=/dev/disk/by-uuid/330fc02e-323e-4b2f-8b02-f7d85086c2fa
Where=/
Type=ext4
Options=errors=remount-ro
